# applications/kit/cron/crontab
#crontab
0 5 * * * root *schedule/limpieza/Modulero
@reboot root *schedule/inicio
